# polybot.service
[Unit]
Description=Polybot Telegram Bot Service
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/home/ubuntu/PolybotService
ExecStart=/home/ubuntu/PolybotService/venv/bin/python3 /home/ubuntu/PolybotService/polybot/bot.py
Restart=always
EnvironmentFile=/home/ubuntu/PolybotService/.env
Environment=PYTHONUNBUFFERED=1


[Install]
WantedBy=multi-user.target
