# polybot-dev.service
[Unit]
Description=Polybot Dev Service
After=network.target

[Service]
ExecStart=/usr/bin/python3 /home/ubuntu/polybot-dev/bot.py
WorkingDirectory=/home/ubuntu/polybot-dev
Restart=always
User=ubuntu
Environment="ENV=development"

[Install]
WantedBy=multi-user.target
