# polybot-dev.service
[Unit]
Description=Polybot Dev Service
After=network.target

[Service]
ExecStart=/usr/bin/env python3 -m polybot.app
WorkingDirectory=/home/ubuntu/PolybotService
Restart=always
User=ubuntu
Environment="ENV=development"

[Install]
WantedBy=multi-user.target


